# docker/root/etc/services.d/ttyd/run
#!/usr/bin/env bash

set -eu
set -o pipefail


THEME="$(cat /etc/ttyd/theme.json)"

/usr/local/bin/ttyd --port=8080 -t "theme=$THEME" /usr/bin/zsh
