# systemd/00-vscode-update.timer
[Unit]
Description=VsCode Server Update Timer

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
